# server/wds/deb/DEBIAN/prerm
#!/bin/sh
set -e

case "$1" in
    remove|purge)
        # Nothing to do here
    ;;
esac

invoke-rc.d nginx stop || true

